# configs/systemd/vault-dm-crypt-refresh.service
[Unit]
Description=Vault DM-Crypt Authentication Refresh
Documentation=https://github.com/axonops/vault-dm-crypt
After=network-online.target
Wants=network-online.target
DefaultDependencies=false

[Service]
Type=oneshot
User=root
Group=root
# Run the refresh-auth command (default: refresh if expiring within 30 minutes and update config)
ExecStart=/usr/bin/vault-dm-crypt refresh-auth
# Set timeout for the operation
TimeoutSec=300
# Log output to journal
StandardOutput=journal
StandardError=journal
# Set nice level for background operation
Nice=10
# Restrict privileges
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ReadWritePaths=/etc/vault-dm-crypt
ProtectHome=true
RestrictSUIDSGID=true
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true
RestrictNamespaces=true
LockPersonality=true
RestrictRealtime=true

[Install]
WantedBy=multi-user.target
