# configs/systemd/vault-dm-crypt-refresh.timer
[Unit]
Description=Vault DM-Crypt Authentication Refresh Timer
Documentation=https://github.com/digitalis-io/vault-dm-crypt
Requires=vault-dm-crypt-refresh.service

[Timer]
# Run every 15 minutes
OnCalendar=*:0/15
# Add random delay up to 5 minutes to prevent thundering herd
RandomizedDelaySec=300
# Start immediately if the system was down during a scheduled run
Persistent=true
# Ensure timer survives system clock changes
AccuracySec=1min

[Install]
WantedBy=timers.target
